# crontab
# 防止hourly中的任务不执行
47 * * * * root /etc/cron.hourly/backup-mysql 2>&1
# Don't remove the empty line at the end of this file. It is required to run the cron job
